# misc/supervisor.conf
# Supervisor config entry for Laravel IP Camera.

[program:laravel-ipcamera]
command=/usr/bin/php /var/www/example.com/artisan camera:watch
numprocs=1
autostart=true
autorestart=true
user=www-data
